# supervisord.conf
[supervisord]
nodaemon=true

[program:traffic_counts]
command=./traffic_counts.sh
autostart=true
autorestart=true